# shared/you.softwareof.brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>you.softwareof.brief</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PLUGIN_ROOT__/scripts/morning-brief.sh</string>
    </array>

    <!-- 7:53 AM, Monday–Friday -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>53</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>53</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>53</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>53</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>53</integer></dict>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_PLUGIN_ROOT</key>
        <string>__PLUGIN_ROOT__</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:__HOME__/.claude/local:__HOME__/.local/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__DATA_DIR__/logs/launchd-brief.log</string>

    <key>StandardErrorPath</key>
    <string>__DATA_DIR__/logs/launchd-brief.log</string>

    <key>WorkingDirectory</key>
    <string>__PLUGIN_ROOT__</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
